# hgots-node.conf
# hgots-node
description "Door lock system in nodejs for Herbartgymnasium Oldenburg."
author "Sören Gade"

start on runlevel [2345]
stop on runlevel [016]

cd /home/pi/hgots-node/
exec /usr/local/bin/node /home/pi/hgots-node/src/app
respawn
